# supervisord.conf
[supervisord]

nodaemon=true

[program:gitlab-runsv]

command=rails server -p 80
